# compose/production/django/crontab
0 0 * * *	root	 /usr/src/osm-extract/osm_update.sh >> /var/log/geonode-task.log 2>&1

